# hmsdk/sysroots/aarch64-euler-elf_all_in_one/m536agrp/m536a/etc/inittab
#
# inittab       This file describes how the INIT process should set up
#		the system during the initial booting process.

# System initialization starts here

::sysinit:/etc/hm.local
::respawn:-/bin/bash -l
::shutdown:/etc/rc.shutdown
