# hmsdk/sysroots/aarch64-euler-elf_all_in_one/m536agrp/m536a/etc/rc.local
#!/bin/bash
echo "This is default rc.local script"

# running local startup scripts
echo "Running local startup scripts."
for f in /etc/rc.start/* ; do
	$f
done
exit 0
